# systemd/familiar-workspace.service
[Unit]
Description=Familiar Workspace (UI + reverse proxy)
After=network.target familiar-gateway.service
Wants=familiar-gateway.service

[Service]
Type=simple
User=youruser
WorkingDirectory=/home/youruser/repos/familiar/familiar-workspace
ExecStart=/home/youruser/repos/familiar/familiar-workspace/familiar-workspace --config /home/youruser/repos/familiar/familiar-workspace/workspace.toml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
